# root/etc/services.d/ac2mqtt/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python3 /app/ac2mqtt/monitor.py \
	-c /config/config.yml